# com.calvin.photosync.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.calvin.photosync.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/kingfisher/.openclaw/workspace/calvin_daemon_enhanced.py</string>
        <string>--interval</string>
        <string>15</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/kingfisher/.calvin_daemon.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/kingfisher/.calvin_daemon.err.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/kingfisher/.openclaw/workspace</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>5</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>